# root/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
Type=
Type=simple
ExecStart=
ExecStart=-/sbin/agetty --skip-login --nonewline --noissue --autologin maki40681 %I $TERM
